# docker/supervisor/supervisord.prod.conf
# Production Supervisor Configuration for FinFinity Portal

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
logfile_maxbytes=50MB
logfile_backups=10

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=100
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
user=root
killasgroup=true
stopasgroup=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=200
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
user=root
killasgroup=true
stopasgroup=true

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=90 --memory=256
autostart=true
autorestart=true
priority=300
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
startsecs=10
startretries=3
killasgroup=true
stopasgroup=true

[program:laravel-schedule]
command=bash -c "while [ true ]; do (php /var/www/html/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
priority=400
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/schedule.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
startsecs=5
startretries=3
killasgroup=true
stopasgroup=true

[program:laravel-websocket]
command=php /var/www/html/artisan websockets:serve --host=0.0.0.0 --port=6001
autostart=false
autorestart=true
priority=500
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/websocket.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
startsecs=5
startretries=3
killasgroup=true
stopasgroup=true

# Log rotation (if crond is available)
[program:logrotate]
command=bash -c "while [ true ]; do logrotate -f /etc/logrotate.conf; sleep 86400; done"
autostart=false
autorestart=true
priority=900
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/logrotate.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface 